# django/deploy/cron.daily/backup
#!/bin/sh

cd /code && /usr/local/bin/inv db.backup borg.backup > /var/log/comses-backup.log 2>&1
