# django/deploy/cron.daily/backup
#!/bin/sh
#
# Run backups over the model library filesystem and db

export DJANGO_SETTINGS_MODULE="core.settings.production"
inv -r /code db.backup borg.backup >> /shared/logs/comses-backup.log 2>&1
