# django/deploy/cron.daily/backup
#!/bin/sh

/usr/local/bin/inv -r /code borg.backup > /var/log/comses-backup.log 2>&1
