# django/deploy/cron.daily/doi_sync_all_metadata
#!/bin/sh

./manage.py doi_sync_metadata --no-interactive --no-dry-run
