# django/deploy/cron.hourly/update_elasticsearch_index
#!/bin/sh

export DJANGO_SETTINGS_MODULE=${DJANGO_SETTINGS_MODULE:-"core.settings.production"}
/code/manage.py update_index
