# django/deploy/cron.weekly/clamav
#!/bin/sh

LOGFILE=/shared/logs/clamav.log-$(date +%F)

freshclam;
clamscan -l ${LOGFILE} -i -r --quiet /shared/library /shared/static
