# django/deploy/cron.weekly/cleanup
#!/bin/sh

export DJANGO_SETTINGS_MODULE="core.settings.production"
./manage.py clearsessions
