# django/deploy/cron.weekly/vacuum
#!/bin/sh

# perform a full vacuum + analyze db and log the results
vacuumdb -f -a -h db -U comsesnet -w -v --analyze > /shared/logs/vacuum.log-$(date +%F) 2>&1
